# tuned.service
[Unit]
Description=Dynamic System Tuning Daemon
After=syslog.target

[Service]
Type=forking
PIDFile=/var/run/tuned/tuned.pid
ExecStartPre=/bin/sh -c "uname -p | grep -qv '^s390x\?$'"
ExecStart=/usr/sbin/tuned -d -c /etc/tuned.conf

[Install]
WantedBy=multi-user.target
