# udev/81-tuned-add-disk.rules
# restart ktune daemon if a new block device is added
ACTION=="add", SUBSYSTEM=="block", ENV{DEVTYPE}=="disk", RUN+="tuned-request-ktune-restart"
