# contrib/upstart/tuned.conf
description "Dynamic system tuning daemon"

start on runlevel [2345]
stop on runlevel [!2345]

exec tuned --log --pid --daemon

post-stop exec rm -f /var/run/tuned/tuned.pid
